# oneclick-kra-sftp-sync.service
[Unit]
Description=Oneclick KRA SFTP Sync Service
After=network.target

[Service]
Type=oneshot
User=your-user
WorkingDirectory=/path/to/oneclick-kra-sftp-sync
ExecStart=/path/to/oneclick-kra-sftp-sync/run.sh run --no-log
EnvironmentFile=/path/to/oneclick-kra-sftp-sync/.env

[Install]
WantedBy=multi-user.target
